# src/devsn2mac.service.in
[Unit]
Description=Configure eth0 MAC address from DEVSN

[Service]
Type=oneshot
ExecStart=@sbindir@/devsn2mac

[Install]
WantedBy=network.target
